# Aws/dht11/temp.service
[Unit]
Description=Temp Service
[Service]
Type=simple
ExecStart=/usr/local/bin/temp.py
[Install]
WantedBy=multi-user.target
