# Linux/Projects/dante_socks/rpmbuild/BUILD/sockd-1.4.3/sockd.service
[Unit]
Description=Sockd Service
[Service]
Type=normal
ExecStart=/usr/local/sbin/sockd
[Install]
WantedBy=multi-user.target
